# deploy/unitree-g1-voice.service
[Unit]
Description=Unitree G1 Voice Interaction Controller
After=network.target pulseaudio.service
Wants=pulseaudio.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/unitree/bk-main
ExecStart=/bin/bash -c 'sleep 5 && /home/unitree/bk-main/scripts/setup_audio.sh && PYTHONPATH=$HOME/.local/lib/python3.8/site-packages:$PYTHONPATH /usr/bin/python3 /home/unitree/bk-main/VoiceInteraction/multimodal_interaction.py'
Restart=on-failure
RestartSec=15

# Python 环境变量
Environment="PYTHONPATH=/home/unitree/.local/lib/python3.8/site-packages:/home/unitree/.local/lib/python3.8/dist-packages"

# PulseAudio 环境变量（使用 root 用户的 PulseAudio）
Environment="PULSE_SERVER=unix:/run/pulse/native"
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket"

# 硬件权限
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/unitree/bk-main

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unitree-g1-voice

[Install]
WantedBy=multi-user.target
