# deploy/unitree-g1-voice.service
[Unit]
Description=Unitree G1 Voice Interaction Controller
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=unitree
Group=unitree
WorkingDirectory=/home/unitree/bk-main
ExecStart=/usr/bin/python3 /home/unitree/bk-main/VoiceInteraction/multimodal_interaction.py
Restart=on-failure
RestartSec=10

# 硬件权限
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/unitree/bk-main

# 日志
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
