# deploy/unitree-g1-voice.service
[Unit]
Description=Unitree G1 Voice Interaction Controller
After=network.target

[Service]
Type=simple
User=unitree
Group=unitree
WorkingDirectory=/home/unitree/bk-main
ExecStart=/usr/bin/python3 /home/unitree/bk-main/VoiceInteraction/multimodal_interaction.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
